# systemd/smartdns.service.in
[Unit]
Description=SmartDNS Server
After=network.target 
StartLimitBurst=0
StartLimitIntervalSec=60

[Service]
Type=forking
PIDFile=@RUNSTATEDIR@/smartdns.pid
EnvironmentFile=@SYSCONFDIR@/default/smartdns
ExecStart=@SBINDIR@/smartdns -p @RUNSTATEDIR@/smartdns.pid $SMART_DNS_OPTS 
KillMode=process
Restart=always
RestartSec=2
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
Alias=smartdns.service
