# package/debian/DEBIAN/prerm
#!/bin/sh

systemctl stop smartdns
systemctl disable smartdns
